# data/supervisord.conf.template
[supervisord]
nodaemon=true

[program:textExtractor]
directory=/usr/data
command=/bin/bash -c ""


[program:pdfPreview]
directory=/usr/data
command=/bin/bash -c ""